# mqtt-history-injector/rootfs/etc/services.d/mqtt-history-injector/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "MQTT History Injector stopped"