# mqtt-history-injector/rootfs/etc/services.d/mqtt-history-injector/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting MQTT History Injector..."

# Запуск основного скрипта
exec python3 /app/run.py