# mqtt-history-injector/rootfs/etc/services.d/mqtt-history-injector/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Запуск MQTT History Injector..."

# Установка Python-зависимостей
pip3 install --no-cache-dir paho-mqtt requests

# Запуск основного скрипта
python3 /usr/bin/run.py